# ymir/skeleton/puppet/modules/ymir/templates/example.conf
# This file is autogenerated by Puppet. Manual changes may be overwritten!
#
# This file illustrates a simple supervisord task definition
#
[program:task]
command=ps aux
numprocs=1
numprocs_start=0
priority=999
autostart=true
autorestart=unexpected
startsecs=1
startretries=10
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/opt/ymir/logs/task.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile=/opt/ymir/logs/task.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=10
environment=HOME=/
